# .config/systemd/user/borgmatic-check.timer
[Unit]
Description=Run borgmatic check

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
