# .config/autostart/xss-lock.desktop
[Desktop Entry]
Type=Application
Name=XSS Lock
Exec=xss-lock.sh
OnlyShowIn=GNOME;Unity;
NoDisplay=true
